# meta-msm/recipes-services/qt-led-service/qt-led.service
[uint]
Description="qt-led service"

[service]
ExecStart=/usr/bin/qt-led
Type=simple

[Install]
WantedBy=multi-user.targetz


